# deploy/cam/rpv-cam.service
[Unit]
Description=RPV Camera Node
After=network.target
Wants=network.target

[Service]
ExecStartPre=/usr/local/bin/rpv-net-setup-pre.sh
ExecStart=/usr/local/bin/rpv-cam
ExecStopPost=/usr/local/bin/rpv-net-teardown.sh
Restart=on-failure
RestartSec=5
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
